# infra/alerts/agency-os-skill-pr-staleness-monitor.timer
[Unit]
Description=Run skill-PR staleness breach monitor every 15 min (KEI-11 Outcome 2)

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
