# infra/alerts/agency-os-alert-pipeline-failure.timer
[Unit]
Description=Run pipeline failure alert every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
