# infra/alerts/agency-os-alert-vendor-budget.timer
[Unit]
Description=Run vendor budget threshold alert hourly on the hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
