# infra/alerts/agency-os-alert-daily-digest.service
[Unit]
Description=Agency OS — last-24h activity digest to Telegram group

[Service]
Type=oneshot
ExecStart=/home/elliotbot/clawd/Agency_OS/.venv/bin/python /home/elliotbot/clawd/Agency_OS/scripts/alerts/daily_digest.py
WorkingDirectory=/home/elliotbot/clawd/Agency_OS
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
