# infra/systemd/completion-sync-worker.service
[Unit]
Description=Agency OS — Completion Sync Worker (KEI-74)
Documentation=https://linear.app/keiracom/issue/KEI-74
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/elliotbot/clawd/Agency_OS/.venv/bin/python3 /home/elliotbot/clawd/Agency_OS/scripts/orchestrator/completion_sync_worker.py
WorkingDirectory=/home/elliotbot/clawd/Agency_OS
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/elliotbot/clawd/logs/completion-sync-worker.log
StandardError=append:/home/elliotbot/clawd/logs/completion-sync-worker.log
MemoryMax=256M
MemoryHigh=200M

[Install]
WantedBy=default.target
