# infra/systemd/governance-freshness-probe.timer
[Unit]
Description=Agency OS — Governance freshness SLO cadence (1h, KEI Agency_OS-cd36)

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=120
Persistent=true
Unit=governance-freshness-probe.service

[Install]
WantedBy=timers.target
