# infra/systemd/memory-core-fact-probe.timer
[Unit]
Description=Agency OS — Memory core-fact drift probe cadence (6h, Agency_OS-zbvs)

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h
RandomizedDelaySec=300
Persistent=true
Unit=memory-core-fact-probe.service

[Install]
WantedBy=timers.target
