# infra/systemd/agents/weaviate-snapshot.service.d/00-onfailure.conf
[Unit]
OnFailure=keiracom-ops-failure-alert@%n.service
