# config/systemd/user/max-inbox-watcher.service
[Unit]
Description=MAX Inbox Watcher (file inbox → max tmux injection)
After=network.target

[Service]
Type=simple
ExecStart=/home/elliotbot/clawd/scripts/max_inbox_watcher.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
