# systemd/agency-cost-rollup.timer
[Unit]
Description=Run agency daily 3-provider cost rollup at 23:55 AEST (13:55 UTC)

[Timer]
OnCalendar=*-*-* 13:55:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
