# systemd/kei45-idle-daemon.timer
[Unit]
Description=KEI-45 idle-daemon timer — every 15 minutes
Documentation=https://linear.app/keiracom/issue/KEI-45

[Timer]
OnBootSec=120
OnUnitActiveSec=900
AccuracySec=30
Persistent=true
Unit=kei45-idle-daemon.service

[Install]
WantedBy=timers.target
