# infra/systemd/agents/agent-memories-indexer.service
[Unit]
Description=Agency OS — public.agent_memories → Weaviate AgentMemories indexer (Agency_OS-lsyd)
Documentation=https://github.com/Keiracom/Agency_OS/pull/1046
After=network-online.target weaviate.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/elliotbot/clawd/Agency_OS
EnvironmentFile=-/home/elliotbot/.config/agency-os/.env
Environment=PYTHONPATH=/home/elliotbot/clawd/Agency_OS/scripts/orchestrator
ExecStart=/home/elliotbot/clawd/Agency_OS/.venv/bin/python3 /home/elliotbot/clawd/Agency_OS/scripts/orchestrator/agent_memories_indexer.py
Restart=on-failure
RestartSec=15
StandardOutput=append:/home/elliotbot/clawd/logs/agent-memories-indexer.log
StandardError=append:/home/elliotbot/clawd/logs/agent-memories-indexer.log
# KEI-44 cgroup pattern — small worker footprint (Supabase + JSON only).
MemoryMax=256M
MemoryHigh=200M

[Install]
WantedBy=default.target
