# infra/systemd/agents/aiden-agent.service
[Unit]
Description=Keiracom aiden agent session (KEI-94 keep-alive)
Documentation=https://linear.app/keiracom/issue/KEI-94
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
Environment="CALLSIGN=aiden"
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
WorkingDirectory=/home/elliotbot/clawd/Agency_OS-aiden
ExecStart=/home/elliotbot/clawd/Agency_OS/scripts/agent_keepalive.sh aiden aiden /home/elliotbot/clawd/Agency_OS-aiden
ExecStartPost=/home/elliotbot/clawd/Agency_OS/scripts/agent_online_notify.sh AIDEN 15 C0B3QB0K1GQ
Restart=always
RestartSec=10
StandardOutput=append:/home/elliotbot/clawd/logs/aiden-agent.log
StandardError=append:/home/elliotbot/clawd/logs/aiden-agent.log

[Install]
WantedBy=default.target
