# infra/systemd/agents/aiden-agent.service
[Unit]
Description=Keiracom aiden agent session
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="CALLSIGN=aiden"
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
ExecStart=/bin/bash -c 'tmux has-session -t aiden 2>/dev/null || tmux new-session -d -s aiden "claude"'
ExecStartPost=/home/elliotbot/clawd/Agency_OS/scripts/agent_online_notify.sh AIDEN 15

[Install]
WantedBy=default.target
