# infra/systemd/agents/heartbeat-writer@.service
[Unit]
Description=Agency OS — KEI-105 heartbeat writer for callsign %i
Documentation=https://linear.app/keiracom/issue/KEI-105
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/elliotbot/clawd/Agency_OS
EnvironmentFile=-/home/elliotbot/.config/agency-os/.env
Environment=CALLSIGN=%i
ExecStart=/home/elliotbot/clawd/Agency_OS/.venv/bin/python3 /home/elliotbot/clawd/Agency_OS/scripts/orchestrator/heartbeat_writer.py
StandardOutput=append:/home/elliotbot/clawd/logs/heartbeat-writer-%i.log
StandardError=append:/home/elliotbot/clawd/logs/heartbeat-writer-%i.log

[Install]
WantedBy=default.target
