# infra/systemd/agents/max-agent.service
[Unit]
Description=Keiracom max agent session
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="CALLSIGN=max"
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
ExecStart=/bin/bash -c 'tmux has-session -t maxbot 2>/dev/null || tmux new-session -d -s maxbot "claude"'
ExecStartPost=/home/elliotbot/clawd/Agency_OS/scripts/agent_online_notify.sh MAX 15 C0B3QB0K1GQ

[Install]
WantedBy=default.target
