# infra/systemd/agents/orion-agent.service
[Unit]
Description=Keiracom orion agent session
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="CALLSIGN=orion"
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
ExecStart=/bin/bash -c 'tmux has-session -t orion 2>/dev/null || tmux new-session -d -s orion "claude"'
ExecStartPost=/home/elliotbot/clawd/Agency_OS/scripts/agent_online_notify.sh ORION 15 C0B3QB0K1GQ

[Install]
WantedBy=default.target
