# infra/systemd/agents/orion-agent.service
[Unit]
Description=Keiracom orion agent session (KEI-94 keep-alive)
Documentation=https://linear.app/keiracom/issue/KEI-94
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
Environment="CALLSIGN=orion"
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
WorkingDirectory=/home/elliotbot/clawd/Agency_OS-orion
ExecStart=/home/elliotbot/clawd/Agency_OS/scripts/agent_keepalive.sh orion orion /home/elliotbot/clawd/Agency_OS-orion
ExecStartPost=/home/elliotbot/clawd/Agency_OS/scripts/agent_online_notify.sh ORION 15 C0B3QB0K1GQ
Restart=always
RestartSec=10
StandardOutput=append:/home/elliotbot/clawd/logs/orion-agent.log
StandardError=append:/home/elliotbot/clawd/logs/orion-agent.log

[Install]
WantedBy=default.target
