# infra/systemd/agents/orion-agent.service
[Unit]
Description=Agency OS Agent — ORION (callsign=orion, tmux=orion)
Documentation=https://linear.app/keiracom/issue/KEI-43
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=/home/elliotbot/clawd/Agency_OS-orion
ExecStart=/home/elliotbot/clawd/Agency_OS-orion/scripts/systemd_agent_supervisor.sh orion orion /home/elliotbot/clawd/Agency_OS-orion
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/elliotbot/clawd/logs/orion-agent.log
StandardError=append:/home/elliotbot/clawd/logs/orion-agent.log

[Install]
WantedBy=default.target
