# infra/systemd/agents/scout-agent.service
[Unit]
Description=Agency OS Agent — SCOUT (callsign=scout, tmux=scout)
Documentation=https://linear.app/keiracom/issue/KEI-43
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=/home/elliotbot/clawd/Agency_OS-scout
ExecStart=/home/elliotbot/clawd/Agency_OS-scout/scripts/systemd_agent_supervisor.sh scout scout /home/elliotbot/clawd/Agency_OS-scout
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/elliotbot/clawd/logs/scout-agent.log
StandardError=append:/home/elliotbot/clawd/logs/scout-agent.log

[Install]
WantedBy=default.target
