# infra/systemd/agents/weaviate-snapshot-r2.service
[Unit]
Description=Daily Weaviate snapshot to Cloudflare R2 (KEI-242)
Documentation=https://linear.app/keiracom/issue/KEI-242
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/home/elliotbot/.config/agency-os/.env
WorkingDirectory=/home/elliotbot/clawd/Agency_OS
ExecStart=/home/elliotbot/clawd/venv/bin/python3 -m src.keiracom_system.backup.weaviate_snapshot
StandardOutput=append:/home/elliotbot/clawd/logs/weaviate-snapshot-r2.log
StandardError=append:/home/elliotbot/clawd/logs/weaviate-snapshot-r2.log

[Install]
WantedBy=default.target
