# infra/systemd/agents/weaviate.service
[Unit]
Description=Weaviate vector database (collective intelligence foundation, KEI-48)
Documentation=https://linear.app/keiracom/issue/KEI-48
After=network-online.target
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=/home/elliotbot/clawd/Agency_OS
EnvironmentFile=-/home/elliotbot/.config/agency-os/.env
ExecStart=/home/elliotbot/clawd/Agency_OS/scripts/orchestrator/weaviate_capped.sh
Restart=on-failure
RestartSec=15
StandardOutput=append:/home/elliotbot/clawd/logs/weaviate.log
StandardError=append:/home/elliotbot/clawd/logs/weaviate.log

[Install]
WantedBy=default.target
